# wendy/wendy.service
[Unit]
Description=wendy Discord Bot
After=network.target

[Service]
Type=simple
User=snerloc
WorkingDirectory=/home/snerloc/discord-bots/wendy
Environment=PERSONA_FILE=/home/snerloc/discord-bots/wendy/personas/wendy.txt
Environment=LLM_MODEL=llama-3.2-3b-uncensored
Environment=ALLOWED_CHANNELS=1503348313397923871,1503348313397923872,1504438057103790120
Environment=ALLOW_BOT_MESSAGES=true
Environment=BOT_NAME=wendy
ExecStart=/usr/bin/node /home/snerloc/discord-bots/shared/bot.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
