# llama-server.service
[Unit]
Description=llama-server LLM Inference
After=network.target

[Service]
Type=simple
User=snerloc
ExecStart=/opt/llama.cpp/build/bin/llama-server \
  --model /opt/models/llama-3.2-3b-uncensored.gguf \
  --port 11434 \
  --threads 8 \
  --parallel 2 \
  --cont-batching \
  --ctx-size 8192 \
  --temp 0.7
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
