# systemd-units/bot4.service
[Unit]
Description=Lyla Discord Bot
After=network.target

[Service]
Type=simple
User=snerloc
WorkingDirectory=/home/snerloc/discord-bots/bot4
ExecStart=/usr/bin/node bot.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
