# systemd-units/bot7.service
[Unit]
Description=Sophie Discord Bot
After=network.target

[Service]
Type=simple
User=snerloc
WorkingDirectory=/home/snerloc/discord-bots/bot7
ExecStart=/usr/bin/node /home/snerloc/discord-bots/shared/bot-runtime.js --persona bot7
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
