# systemd-units/wendy.service
[Unit]
Description=Wendy Discord Bot
After=network.target

[Service]
Type=simple
User=snerloc
WorkingDirectory=/home/snerloc/discord-bots/wendy
Environment=PERSONA_FILE=/home/snerloc/discord-bots/wendy/personas/wendy.txt
Environment=LLM_MODEL=dolphin-2.8-mistral-7b-v02
Environment=ALLOWED_CHANNELS=1503348313397923871,1503348313397923872
ExecStart=/usr/bin/node /home/snerloc/discord-bots/shared/bot.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
